# systemd/claude-vector.slice
[Unit]
Description=Claude tool farm — vector DB slice

[Slice]
MemoryAccounting=yes
IOAccounting=yes
TasksAccounting=yes

# Medium-repo defaults — tune per workload
MemoryHigh=8G
MemoryMax=16G
MemorySwapMax=4G
IOWeight=80
